# deploy/agentibridge.service
[Unit]
Description=AgentiBridge MCP Server (native)
After=network.target agentibridge-db.service
Wants=agentibridge-db.service

[Service]
Type=simple
EnvironmentFile=-%h/.agentibridge/docker.env
Environment=AGENTIBRIDGE_TRANSPORT=sse
ExecStart=%h/.venv/bin/python3 -m agentibridge
WorkingDirectory=%h
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
