# deploy/agentibridge-native.service
[Unit]
Description=AgentiBridge MCP Server (Native)
After=network.target redis-server.service

[Service]
Type=simple
EnvironmentFile=%h/.config/agentibridge/env
ExecStart=%h/.local/bin/python3 -m agentibridge
WorkingDirectory=%h
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
