# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:redis]
command=redis-server --appendonly yes --dir /data
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:agentibridge]
command=python3 -m agentibridge
directory=/app
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=
    REDIS_URL="redis://127.0.0.1:6379/0",
    AGENTIBRIDGE_TRANSPORT="sse",
    AGENTIBRIDGE_HOST="0.0.0.0",
    AGENTIBRIDGE_PORT="8100",
    PYTHONPATH="/app"
